# services/meeting-service/supervisord.conf
# supervisord config for the Citadel meeting media stack
# (aceteam-ai/citadel-cli#514). One container, three supervised long-lived
# processes: the virtual display, the audio server, and the meetingd control API.
# Chromium and ffmpeg are NOT supervised here -- they are per-session children
# that meetingd spawns and reaps.
#
# supervisord runs as root (PID 1 under dumb-init) so it can write the log fds and
# manage the child processes, but every program drops to the non-root `bot` user.
# Priorities enforce boot ordering: Xvfb and PulseAudio come up before meetingd so
# meetingd's first /health canary has a display and a pulse server to probe.
[supervisord]
nodaemon=true
user=root
loglevel=info
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_XVFB_RESOLUTION)s -nolisten tcp
user=bot
priority=10
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pulseaudio]
# -n -F loads ONLY our minimal citadel.pa (no system autodetect); exit-idle-time
# -1 + disallow-exit keep the server up for the container's whole life.
command=/usr/bin/pulseaudio -n -F /etc/pulse/citadel.pa --daemonize=no --exit-idle-time=-1 --disallow-exit --log-target=stderr
user=bot
priority=20
autorestart=true
startsecs=2
environment=HOME="/home/bot",XDG_RUNTIME_DIR="/run/user/10001"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cdp-forward]
# Chromium binds its DevTools socket to 127.0.0.1 only (it refuses a non-loopback
# bind), which a docker port publish cannot reach. socat bridges the
# container-external CDP port to Chromium's loopback port so the published
# 127.0.0.1:HOST:CDP mapping works. reuseaddr + fork so it survives Chromium
# restarts between meetings; it simply refuses connections while no session is up.
command=/usr/bin/socat TCP-LISTEN:%(ENV_MEETING_CDP_PORT)s,fork,reuseaddr,bind=0.0.0.0 TCP:127.0.0.1:%(ENV_MEETING_CDP_INTERNAL_PORT)s
user=bot
priority=25
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:meetingd]
command=/opt/venv/bin/uvicorn meetingd:app --host 0.0.0.0 --port %(ENV_MEETINGD_PORT)s
directory=/app
user=bot
priority=30
autorestart=true
startsecs=3
# meetingd binds 0.0.0.0 INSIDE the container: a docker port publish maps the host
# port to the container's eth0, not its loopback, so a 127.0.0.1 bind here would be
# unreachable through the published port. The compose publish
# (127.0.0.1:HOST:CONTAINER) is what restricts it to host loopback.
environment=HOME="/home/bot",XDG_RUNTIME_DIR="/run/user/10001",PULSE_SERVER="unix:/run/user/10001/pulse/native",DISPLAY="%(ENV_DISPLAY)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
